# systemd/build-server-runner.service
[Unit]
Description=Build server runner — spawns ephemeral GitHub Actions runner containers
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/build-server
EnvironmentFile=/opt/build-server/.env
ExecStart=/opt/build-server/venv/bin/uvicorn runner_pool.runner:app --host 127.0.0.1 --port 3001
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
